# cloud-init/usr/lib/systemd/system/serial-getty@ttyS0.service.d/10-autologin.conf
# Enable automatic root login on cloud
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty --autologin root --noclear %I $TERM
